# packaging/DEBIAN/postinst
#!/bin/bash
set -e

# Step 1: Set executable permissions on the Python scripts
chmod +x /usr/lib/jiopc-hibernate/apply_handlers.py || true
chmod +x /usr/lib/jiopc-hibernate/restore_service.py || true

# Step 2: Create the hibernate state directory under the JioPC user's home
mkdir -p /home/user/.local/share/jiopc/hibernate || true
chown -R user:user /home/user/.local/share/jiopc || true

# Step 3: Set correct permissions on the installed config/autostart files
chmod 644 /usr/lib/jiopc-hibernate/handlers/handlers.yaml || true
chmod 644 /etc/xdg/autostart/jiopc-restore.desktop || true
# Fallback in case it was installed to usr/share/autostart instead
chmod 644 /usr/share/autostart/jiopc-restore.desktop 2>/dev/null || true

# Step 4: Set permissions on the systemd user service unit
chmod 644 /etc/systemd/user/jiopc-save-on-logout.service || true

# Step 5: Enable the service globally so it auto-starts on every future login
# for all users on this machine (no user action ever required).
systemctl --global enable jiopc-save-on-logout.service 2>/dev/null || true

# Step 6: Also activate it immediately for the currently logged-in session
# so dpkg -i takes effect right now — no logout/login cycle needed.
# postinst runs as root, so we must supply the user's runtime environment
# explicitly for systemctl --user to reach the correct D-Bus daemon.
JIOPC_USER="user"
JIOPC_UID=$(id -u "$JIOPC_USER" 2>/dev/null) || true
JIOPC_RUNTIME="/run/user/$JIOPC_UID"

if [ -n "$JIOPC_UID" ] && [ -d "$JIOPC_RUNTIME" ]; then
    runuser -u "$JIOPC_USER" -- \
        env XDG_RUNTIME_DIR="$JIOPC_RUNTIME" \
            DBUS_SESSION_BUS_ADDRESS="unix:path=$JIOPC_RUNTIME/bus" \
        systemctl --user daemon-reload 2>/dev/null || true

    runuser -u "$JIOPC_USER" -- \
        env XDG_RUNTIME_DIR="$JIOPC_RUNTIME" \
            DBUS_SESSION_BUS_ADDRESS="unix:path=$JIOPC_RUNTIME/bus" \
        systemctl --user start jiopc-save-on-logout.service 2>/dev/null || true
fi

# Step 7: Print to stdout
echo "jiopc-session-hibernate installed successfully." || true

exit 0
