# packaging/etc/systemd/user/jiopc-save-on-logout.service
[Unit]
Description=JioPC Session Hibernate - Automatically capture session state on logout
Documentation=https://github.com/jiopc/session-hibernate

[Service]
Type=oneshot
# RemainAfterExit keeps this service "active" for the entire session lifetime.
# systemd will then run ExecStop when the user's session is torn down on logout.
RemainAfterExit=yes
Environment=DISPLAY=:0
ExecStart=/bin/true

# Component A: fires automatically on ANY logout path — LXQt panel button,
# network drop, VM timeout, xss-lock idle — no user action required.
ExecStop=bash -c "cd /usr/lib/jiopc-hibernate && python3 -m save.save_service"

# Match the 10-second capture budget from design.md, with 2s grace margin.
TimeoutStopSec=12

[Install]
# default.target is ALWAYS activated on every user login, regardless of
# desktop environment. More reliable than graphical-session.target which
# LXQt does not fully activate on this system.
WantedBy=default.target
